# dhivehi-transliteration.service
# Systemd service for Dhivehi Transliteration
#
# Setup:
# 1. Edit paths below (User, WorkingDirectory, Environment, ExecStart)
# 2. sudo cp dhivehi-transliteration.service /etc/systemd/system/
# 3. sudo systemctl daemon-reload
# 4. sudo systemctl enable dhivehi-transliteration
# 5. sudo systemctl start dhivehi-transliteration

[Unit]
Description=Dhivehi Transliteration Web Service
After=network.target

[Service]
Type=notify
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=/path/to/div-transliteration
Environment="PATH=/path/to/div-transliteration/venv/bin"
ExecStart=/path/to/div-transliteration/venv/bin/gunicorn -c gunicorn.conf.py app:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
